# upstart/reddit_startup.conf
start on runlevel [2345]
stop on runlevel [!2345]

setuid jerry
setgid jerry

chdir /home/jerry/Reddit_Craigslist_Imgur_Bot/
exec python app.py